# scripts/proton/crontab
ZPROFILE=/home/pi/.zprofile
CRON_SCRIPTS=/home/pi/dotfiles/scripts/proton

0 0 * * * . "$ZPROFILE"; $CRON_SCRIPTS/backup.sh
