# scripts/electron/crontab
ZPROFILE=/home/henrik/.zprofile
CRON_SCRIPTS=/home/henrik/dotfiles/scripts/electron

0 0 * * * . "$ZPROFILE"; $CRON_SCRIPTS/backup.sh
